# ops/systemd/monolith-web.service
[Unit]
Description=Mon-olith Web (Next.js)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/srv/Monolith/apps/web
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=NEXT_TELEMETRY_DISABLED=1
Environment=NEXT_PUBLIC_API_BASE_URL=https://monolith-labs.xyz
ExecStart=/usr/bin/env NODE_ENV=production /usr/bin/npm run start
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
